# linux/psychocity-rhel/boot.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
After=blockdev@dev-disk-by\x2duuid-24a055ee\x2ddd70\x2d4203\x2daa6a\x2d8d295be6f4ec.target

[Mount]
What=/dev/disk/by-uuid/24a055ee-dd70-4203-aa6a-8d295be6f4ec
Where=/boot
Type=xfs
