# linux/psychocity-rhel/sshd-unix-local.socket
# Automatically generated by systemd-ssh-generator

[Unit]
Description=OpenSSH Server Socket (systemd-ssh-generator, AF_UNIX Local)
Documentation=man:systemd-ssh-generator(8)

[Socket]
ListenStream=/run/ssh-unix-local/socket
Accept=yes
PollLimitIntervalSec=30s
PollLimitBurst=50
